# duckdns.service
[Unit]
Description=DuckDNS update job

[Service]
Type=oneshot
ExecStart=/usr/local/bin/duckdns
